# systemd/kronos.socket
[Unit]
Description=Kronos rest api

[Socket]
ListenStream=10000

[Install]
WantedBy=sockets.target
