# systemd/kronos-cluster-node.socket
[Socket]
ListenStream=/run/{{name}}/socket
MaxConnections=4

[Install]
RequiredBy={{name}}.service
